# vls@.service.in
[Unit]
Description=Video Looped Server

[Service]
Type=simple
ExecStart=@install_path@/vls --listening_address=0.0.0.0 %I

[Install]
WantedBy=multi-user.target
